# ops/systemd/argos-frontend.service.template
[Unit]
Description=Argos Frontend Service
After=network.target argos-backend.service

[Service]
User=argos
Group=argos
WorkingDirectory=/opt/argos/frontend
EnvironmentFile=/etc/argos/argos.env
# Serve built static files (build with: pnpm run build)
ExecStart=/usr/bin/python3 -m http.server 5173 --directory dist
Restart=always

[Install]
WantedBy=multi-user.target
